# gpio_shutdown_daemon/gpio_shutdown_daemon.conf
[program:gpio_shutdown_daemon]
command=sudo /home/pi/gpio_shutdown_daemon &
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/gpio_shutdown_daemon.log
user=pi
